# glugbot.timer
[Unit]
Description=glugbot

[Timer]
OnCalendar=daily
OnStartupSec=60m
Persistent=true

[Install]
WantedBy=timers.target
