# glugbot.timer
[Unit]
Description=glugbot

[Timer]
OnCalendar=*-*-* 11:00:00
OnStartupSec=60m
Persistent=true

[Install]
WantedBy=timers.target
